# data/openrc/gnome-shell
#!/sbin/openrc-run

description="GNOME Shell"
XDG_CURRENT_DESKTOP=GNOME
export XDG_CURRENT_DESKTOP
command="/usr/bin/gnome-shell"
command_args="--wayland"
command_background="true"
pidfile="${XDG_RUNTIME_DIR}/gnome-shell.pid"
