# data/openrc/gnome-shell-wayland.gnome-login
#!/sbin/openrc-run

supervisor=supervise-daemon
description="GNOME Shell (login)"
XDG_CURRENT_DESKTOP=GNOME
export XDG_CURRENT_DESKTOP
command="/usr/bin/gnome-shell"
command_args="--wayland --mode=gdm"
command_background="true"
pidfile="${XDG_RUNTIME_DIR}/gnome-shell-gdm.pid"
