# data/openrc/gnome-shell-x11
#!/sbin/openrc-run

XDG_CURRENT_DESKTOP=GNOME
export XDG_CURRENT_DESKTOP
XDG_SESSION_TYPE=x11
export XDG_SESSION_TYPE
# HACK: Need better openrc env handling
DBUS_SESSION_BUS_ADDRESS="unix:path=${XDG_RUNTIME_DIR}/bus"
export DBUS_SESSION_BUS_ADDRESS
XDG_DATA_DIRS="$XDG_DATA_DIRS:/home/$USER/.local/share/flatpak/exports/share:/var/lib/flatpak/exports/share"
export XDG_DATA_DIRS
DISPLAY=:2
export DISPLAY

supervisor=supervise-daemon
description="GNOME Shell (x11)"
command="/usr/bin/gnome-shell"
command_args="--x11"
#command_background="true"
pidfile="${XDG_RUNTIME_DIR}/gnome-shell-x11.pid"
