# data/openrc/gnome-session-dbus
#!/sbin/openrc-run

_session="${RC_SVCNAME#*.}"

supervisor=supervise-daemon
description="GNOME session dbus daemon"
command="/usr/libexec/gnome-session-service"
comand_args='--session="${_session}"'
command_background="true"
pidfile="${XDG_RUNTIME_DIR}/gnome-session-dbus.pid"
