# data/openrc/gnome-session-x11
#!/sbin/openrc-run

description="GNOME Session (x11)"

gs_x11_session="${RC_SVCNAME#*.}"

depend() {
	need dbus
	if [ "${gs_x11_session}" = "gnome-login" ]; then
		need gnome-shell-x11.gnome-login
	else
		need gnome-shell-x11
	fi
	# need "gnome-session.${gs_x11_session}"
	need gnome-settings-daemon-x11
	need "gnome-session-dbus.${gs_x11_session}"
	need gnome-session-monitor
	need gnome-session-init
}
