# data/openrc/gnome-session
#!/sbin/openrc-run

description="GNOME Session (base)"

depend() {
	need gnome-shell
	need gnome-settings-daemon
	need gnome-session-dbus
}
