# data/openrc/gnome-session-wayland
#!/sbin/openrc-run

description="GNOME Session (wayland)"

_session="${RC_SVCNAME#*.}"

depend() {
	need "gnome-session.${_session}"
}
