# data/openrc/gnome-session
#!/sbin/openrc-run

supervisor=supervise-daemon
description="GNOME Session (base)"
output_logger="logger"
error_logger="logger"

gnome_session="${RC_SVCNAME#*.}"

depend() {
	need dbus
	need gnome-settings-daemon
	need "gnome-session-dbus.${gnome_session}"
	need gnome-session-monitor
	need gnome-session-init
}

