# data/openrc/gnome-session-wayland
#!/sbin/openrc-run

description="GNOME Session (wayland)"

_session=$(echo "${RC_SVCNAME}" | sed -e "s|.*\.||g")

depend() {
	need "gnome-session.${_session}"
}
