# data/openrc/gnome-session-wayland
#!/sbin/openrc-run

description="GNOME Session (wayland)"
output_logger="logger"
error_logger="logger"

_session="${RC_SVCNAME#*.}"

depend() {
	if [ "${_session}" = "gnome-login" ]; then
		need gnome-shell-wayland.gnome-login
	else
		need gnome-shell-wayland
	fi
	need "gnome-session.${_session}"
}
