# data/openrc/gnome-session-wayland
#!/sbin/openrc-run

description="GNOME Session (wayland)"

depend() {
	need gnome-session
}
