# data/openrc/gnome-session-init
#!/sbin/openrc-run

supervisor=supervise-daemon
description="Gnome session init signal"
command="/usr/libexec/gnome-session-ctl"
command_args="--signal-init"
pidfile="${XDG_RUNTIME_DIR}/gnome-session-init.pid" #?
notify="socket:ready"
