# data/openrc/gnome-shell
#!/sbin/openrc-run

XDG_CURRENT_DESKTOP=GNOME
export XDG_CURRENT_DESKTOP
# HACK: Need better openrc env handling
DBUS_SESSION_BUS_ADDRESS="unix:path=${XDG_RUNTIME_DIR}/bus"
export DBUS_SESSION_BUS_ADDRESS

supervisor=supervise-daemon
description="GNOME Shell"
command="/usr/bin/gnome-shell"
command_args="--wayland"
#command_background="true"
pidfile="${XDG_RUNTIME_DIR}/gnome-shell.pid"
