# data/openrc/gnome-session-x11
#!/sbin/openrc-run

description="GNOME Session (x11)"

_session="${RC_SVCNAME#*.}"

depend() {
	if [ "${_session}" = "gnome-login" ]; then
		need gnome-shell-x11.gnome-login
	else
		need gnome-shell-x11
	fi
	need "gnome-session.${_session}"
}
