# data/openrc/gnome-session-init
#!/sbin/openrc-run

# HACK: Need better openrc env handling
DBUS_SESSION_BUS_ADDRESS="unix:path=${XDG_RUNTIME_DIR}/bus"
export DBUS_SESSION_BUS_ADDRESS

supervisor=supervise-daemon
description="Gnome session init signal"
command="/usr/libexec/gnome-session-ctl"
command_args="--signal-init"
pidfile="${XDG_RUNTIME_DIR}/gnome-session-init.pid" #?
# notify="socket:ready"
