# data/openrc/gnome-session
#!/sbin/openrc-run

description="GNOME Session (base)"

_session=$(echo "${RC_SVCNAME}" | sed -e "s|.*\.||g")

depend() {
	need dbus
	if [ "${_session}" == "gnome-login" ]; then
		need gnome-shell.gnome-login
	else
		need gnome-shell
	fi
	need gnome-settings-daemon
	need "gnome-session-dbus.${_session}"
}
