# data/openrc/gsd
#!/sbin/openrc-run

description="${RC_SVCNAME} daemon"
command="/usr/libexec/${RC_SVCNAME}"
command_args=""
command_background="true"
pidfile="${XDG_RUNTIME_DIR}/${RC_SVCNAME}.pid"
